# imapfilter.conf
description "imapfilter lua script remotely sorting my mails"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 1 10

exec su -c 'imapfilter' YourUnixName
